# .config/autostart/xfce4_guake.desktop
[Desktop Entry]
Name=Guake
GenericName=Terminal emulator
Comment=Guake
Exec=guake
Terminal=false
Type=Application
Icon=guake
