# .config/autostart/xfce4_slack.desktop
[Desktop Entry]
Name=Slack
GenericName=Slack
Comment=Media rich messaging service
Exec=slack
Terminal=false
Type=Application
Icon=slack
