# home/.config/autostart/xfce4_keeweb.desktop
[Desktop Entry]
Name=KeeWeb
GenericName=Credential storage
Comment=Auth
Exec=KeeWeb
Terminal=false
Type=Application
Icon=keeweb
